# resources/fty-outage.service.in
[Unit]
Description=fty-outage service
After=malamute.service network.target
Requires=malamute.service network.target

# Enable FLA... once, do not keep retrying if EULA is not accepted yet
# There are no assets monitored and metrics produced before FLA anyway
After=fty-license-accepted.target
Wants=fty-license-accepted.target
Requisite=fty-license-accepted.target

PartOf=bios.target

[Service]
Type=simple
User=bios
Restart=always
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/share/bios/etc/default/bios
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/share/bios/etc/default/bios__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/share/fty/etc/default/fty
EnvironmentFile=-@CMAKE_INSTALL_PREFIX@/share/fty/etc/default/fty__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_SYSCONFDIR@/default/bios
EnvironmentFile=-@CMAKE_INSTALL_SYSCONFDIR@/default/bios__%n.conf
EnvironmentFile=-@CMAKE_INSTALL_SYSCONFDIR@/default/fty
EnvironmentFile=-@CMAKE_INSTALL_SYSCONFDIR@/default/fty__%n.conf
Environment="prefix=@CMAKE_INSTALL_PREFIX@"
ExecStart=@CMAKE_INSTALL_PREFIX@/bin/@PROJECT_NAME@ -c @CMAKE_INSTALL_SYSCONFDIR@/@PROJECT_NAME@/@PROJECT_NAME@.cfg

[Install]
WantedBy=bios.target
